# crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
0 20 * * * root /usr/local/bin/github_sync.sh > /proc/1/fd/1 2>&1
