# crontab
0 * * * * /usr/local/bin/github_sync.sh
